# packaging/debian/ovd-administration-console/lucid/ulteo-ovd-administration-console.postinst
#!/bin/sh -e

. /usr/share/debconf/confmodule

CONFDIR=/etc/ulteo/administration_console
CONFFILE=$CONFDIR/config.inc.php
SPOOLDIR=/var/spool/ulteo/administration_console
INSTALLDIR=/usr/share/ulteo/administration_console

############################# APACHE CONFIGURATION #############################

if [ -d /etc/apache2/conf-available ]; then
    A2CONFDIR=/etc/apache2/conf-available
else
    A2CONFDIR=/etc/apache2/conf.d
fi

A2USER=www-data

# Alias admin
if [ ! -e $A2CONFDIR/ovd-administration-console.conf ]; then
    ln -sfT $CONFDIR/apache2.conf $A2CONFDIR/ovd-administration-console.conf
    which a2enconf && a2enconf ovd-administration-console
fi

# Update wsdl
if [ -f $INSTALLDIR/tools/update_wsdl_cache.php ]
then
   echo "Purging wsdl cache files."
   su -s /usr/bin/php $A2USER "$INSTALLDIR/tools/update_wsdl_cache.php"
   [ $? -ne 0 ] && exit 1
fi

# restart apache server
if apache2ctl configtest 2>/dev/null; then
    invoke-rc.d apache2 reload || true
else
    echo << EOF
Your apache configuration is broken!
Correct it and restart apache.
EOF
fi

# create folders and change permissions
mkdir -p $SPOOLDIR
chown $A2USER:$A2USER $SPOOLDIR $CONFFILE

############################### FIRST INSTALL #################################

if [ "$1" = "configure" ] && ([ -z "$2" ] || [ -n "$DEBCONF_RECONFIGURE" ])
then
    db_get ulteo-ovd-administration-console/session_manager_host && SM_HOST=$RET

    # set the configuration
    sed -i "/SESSIONMANAGER_HOST/ s/, *'.*' *);/, '${SM_HOST}');/" $CONFFILE
fi

db_stop

#DEBHELPER#

exit 0
