# packaging/debian/ovd-administration-console/lenny/ulteo-ovd-administration-console.prerm
#!/bin/sh -e

A2CONFDIR=/etc/apache2/conf.d

# remove apache2 config link
rm -f $A2CONFDIR/ovd-administration-console.conf

# reload apache
if apache2ctl configtest 2>/dev/null; then
    invoke-rc.d apache2 reload || true
else
    echo "Your apache2 configuration is broken, so we're not restarting it for you."
fi

#DEBHELPER#

exit 0
