# packaging/debian/ovd-webinterface/lucid/ulteo-ovd-webinterface.postrm
#!/bin/sh -e

A2CONFDIR=/etc/apache2/conf.d

if [ "$1" = "remove" ]; then
    # remove apache2 link
    rm -f /etc/apache2/conf.d/webinterface

    # reload apache
    if apache2ctl configtest 2>/dev/null; then
        invoke-rc.d apache2 reload || true
    else
        echo "Your apache2 configuration is broken, so we're not restarting it for you."
    fi
fi

#DEBHELPER#

exit 0
