# SessionManager/config/sessionmanager.cron.in
# /etc/cron.d/sessionmanager: crontab fragment for ulteo-ovd-sessionmanager

 * * * * * @APACHE_USER@ /usr/bin/php @DATADIR@/ulteo/sessionmanager/crons/minutely.php >/dev/null 2>&1
00 * * * * @APACHE_USER@ /usr/bin/php @DATADIR@/ulteo/sessionmanager/crons/hourly.php >/dev/null 2>&1
